# docker/supervisord.conf
[unix_http_server]
file=/tmp/supervisor.sock
chmod=0700
chown=www-data:www-data

[supervisord]
nodaemon=true
user=www-data
logfile=/dev/null
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[program:php-fpm]
command=php-fpm -F
user=www-data
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:nginx]
command=nginx -g "daemon off;"
user=www-data
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

# ── Queue worker ─────────────────────────────────────────────────────────────
# No explicit connection: works whatever QUEUE_CONNECTION is set to
# (database, redis, ...). --max-time lets supervisor restart it periodically
# to release memory.
[program:queue_worker]
command=php /var/www/html/artisan queue:work --sleep=3 --tries=3 --max-time=3600
user=www-data
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/queue_worker.out.log
stdout_logfile_maxbytes=0
stderr_logfile=/var/log/supervisor/queue_worker.err.log
stderr_logfile_maxbytes=0

# ── Scheduler ────────────────────────────────────────────────────────────────
# Poor-man's cron for containers: runs schedule:run every 60s.
[program:scheduler]
command=sh -c 'while true; do php /var/www/html/artisan schedule:run --no-interaction; sleep 60; done'
user=www-data
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/scheduler.out.log
stdout_logfile_maxbytes=0
stderr_logfile=/var/log/supervisor/scheduler.err.log
stderr_logfile_maxbytes=0
